# nginx-app/installs/configs.sh
#!/sbin/openrc-run

adduser -D -H -u 1000 -s /bin/bash -Gwww-data www-data
